# miscellaneous/wg-quick-oci.target
[Unit]
Description=WireGuard Tunnels via wg-quick-oci
